# Backend/backend.service

[Unit]
Description=Broker Service For Noetic
Requires=rabbitmq-server.service
After=network.target

[Service]
ExecStart=/home/it490/IT-490-002/Backend/db_listener.php
Restart=always
User=it490

[Install]
WantedBy=multi-user.target
