# DMZ/dmz.service
[Unit]
Description=DMZ For External Connections Controlled by Jess
After=network.target

[Service]


[Install]
WantedBy=multi-user.target
